# deploy/systemd/resort-os-backup.service
# deploy/systemd/resort-os-backup.service
#
# Runs scripts/backup_db.sh once. Triggered on a schedule by the paired
# resort-os-backup.timer — not enabled directly.
#
# Install (adjust WorkingDirectory/User to match your server):
#   sudo cp deploy/systemd/resort-os-backup.service deploy/systemd/resort-os-backup.timer /etc/systemd/system/
#   sudo systemctl daemon-reload
#   sudo systemctl enable --now resort-os-backup.timer
#   systemctl list-timers resort-os-backup.timer   # confirm next run time
#   journalctl -u resort-os-backup.service          # check logs after it fires
#
# Offsite sync to S3/Backblaze (wagdy.md T-04, optional) is controlled by
# BACKUP_REMOTE_ENABLED/BACKUP_RCLONE_REMOTE/BACKUP_RCLONE_CONFIG in
# backend/.env (same file DATABASE_URL is read from) — not set here, so this
# unit file doesn't need touching to turn it on/off. If it's enabled and the
# sync fails, this service reports FAILED (non-zero exit) even though the
# local dump still succeeded — check journalctl, the local backup is safe.

[Unit]
Description=Resort OS — PostgreSQL backup
After=docker.service

[Service]
Type=oneshot
User=resortos
WorkingDirectory=/opt/wegosharm/resort-os
ExecStart=/opt/wegosharm/resort-os/scripts/backup_db.sh
Environment=BACKUP_RETENTION_DAYS=14
