# deploy/systemd/resort-os-backup.timer
# deploy/systemd/resort-os-backup.timer
#
# Fires resort-os-backup.service daily at 03:00 server time — chosen to sit
# well clear of the app's own scheduled jobs in app/celery_app.py (night
# audit at 00:01, most reminder jobs at 08:00-10:00), so a slow backup never
# competes with those for DB I/O.
#
# Persistent=true — if the server was off/rebooting at 03:00, the backup
# runs once at next boot instead of silently being skipped for a day.

[Unit]
Description=Daily Resort OS PostgreSQL backup

[Timer]
OnCalendar=*-*-* 03:00:00
Persistent=true
RandomizedDelaySec=300

[Install]
WantedBy=timers.target
